# atsar.cron
#0     *    * * 0-6 root [ -f /var/lock/subsys/atsar ] && //usr/bin/atsa1
#20,40 8-17 * * 1-5 root [ -f /var/lock/subsys/atsar ] && //usr/bin/atsa1
 *     *    * * 0-6 root [ -f /var/lock/subsys/atsar ] && //usr/bin/atsa1
